# camservice
### BEGIN INIT INFO
# Provides:          camservice
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start camservice at boot time
# Description:       Start camservice at boot time.
### END INIT INFO

#!/bin/bash
# /etc/init.d/camservice

### BEGIN INIT INFO
### END INIT INFO


case "$1" in 
    start)
        echo "Starting camservice"
        sudo ./home/pi/ngrok authtoken 46b4Uk6yjGbVxum9E5WSA_5BTdsiUqfCtRUhSy11mP4
        #sudo nohup /home/pi/ngrok http 80 &
        sudo nohup /home/pi/corcovadoStream/autoCam.sh &
        ;;
    stop)
        echo "Stopping camservice"
        sudo killall autoCam.sh
        #sudo killall ngrok
        ;;
    *)
        echo "Usage: /etc/init.d/camservice start|stop"
        exit 1
        ;;
esac

exit 0